# agent/scripts/curato-agent.service
[Unit]
Description=Curato Agent - Museum Display Management
Documentation=https://github.com/sagrkv/curato-app01
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=curato
Group=curato
WorkingDirectory=/opt/curato/agent
ExecStart=/usr/bin/node dist/index.js
Restart=always
RestartSec=10
StartLimitIntervalSec=300
StartLimitBurst=5

# Environment
Environment=NODE_ENV=production
EnvironmentFile=-/opt/curato/agent/.env

# Security hardening
ProtectSystem=strict
ReadWritePaths=/opt/curato/agent /var/log/curato /tmp
PrivateTmp=true
NoNewPrivileges=true
ProtectHome=true
ProtectKernelTunables=true
ProtectControlGroups=true
RestrictSUIDSGID=true

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=curato-agent

[Install]
WantedBy=multi-user.target
