# server/deploy/museumos.service
[Unit]
Description=Museum OS Server
After=network.target postgresql.service mosquitto.service
Wants=postgresql.service

[Service]
Type=simple
User=museumos
Group=museumos
WorkingDirectory=/opt/museumos/server
ExecStart=/usr/bin/node dist/index.js
Restart=always
RestartSec=5
Environment=NODE_ENV=production
EnvironmentFile=/opt/museumos/server/.env

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/museumos/server/storage /opt/museumos/server/backups /opt/museumos/server/logs

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=museumos

[Install]
WantedBy=multi-user.target
